# contrib/systemd/chessbots.service
[Unit]
Description=Chess Bots Runner
After=network.target

[Service]
Type=simple
User=popiol
WorkingDirectory=/home/popiol/chess_bots
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python -m src.agents.runner
# Ensure the engine binary is discoverable for systemd services (includes /usr/games)
Environment=STOCKFISH_PATH=/usr/games/stockfish
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games
# Forward stdout/stderr to the journal so tracebacks are visible via journalctl
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
